# systemd/restic-dl.socket
[Socket]
ListenStream = 9242

[Install]
WantedBy = sockets.target